# voice-agent/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:voice-agent]
command=python agent.py start
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=5
stderr_logfile=/var/log/voice-agent.err.log
stdout_logfile=/var/log/voice-agent.out.log
environment=PYTHONUNBUFFERED="1"

[program:ui-server]
command=python ui_server.py
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=3
stderr_logfile=/var/log/ui-server.err.log
stdout_logfile=/var/log/ui-server.out.log
environment=PYTHONUNBUFFERED="1"
